# rootfs/etc/s6-overlay/s6-rc.d/php7-fpm/run
#!/command/execlineb -P
/usr/sbin/php-fpm7 --force-stderr --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/www.conf


