# rootfs/etc/services.d/horizon/run
#!/bin/sh

if [[ "$ENABLE_LARAVEL_HORIZON" == "1" ]]; then
    php artisan horizon
fi
