# rootfs/etc/services.d/nginx/run
#!/bin/sh
/usr/sbin/nginx -g "daemon off; error_log /dev/stderr info;"
echo "nginx running"
