# rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
/usr/sbin/nginx -g "daemon off; error_log /dev/stderr warn;"
